# scripts/email-assistant.service
[Unit]
Description=Email Assistant Server
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/projects/email-assistant
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10

# 输出日志到systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=email-assistant

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ghm/projects/email-assistant

[Install]
WantedBy=multi-user.target